# watchyourlan/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash

# Wait for WatchYourLAN to be available
bashio::net.wait_for 8840
bashio::log.info "Starting NIGNX..."
exec nginx